# rc.local
#!/bin/bash

/usr/bin/loadkeys /usr/share/keymaps/Caps2Ctrl.map
